# home/dot_config/systemd/user/cnc_relogin.service.tmpl
[Unit]
Description=Logout then login campus network (05:00 daily via timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

Environment=CNC_USER_ID={{ .env.CNC_USER_ID | quote }}
Environment=CNC_PASSWORD={{ .env.CNC_PASSWORD | quote }}
Environment=CNC_SERVICE={{ .env.CNC_SERVICE | quote }}
ExecStartPre=/home/chao/.local/bin/cnc logout
ExecStart=/usr/bin/sleep 2
ExecStartPost=/home/chao/.local/bin/cnc login

StandardOutput=journal
StandardError=journal
