# home/dot_config/systemd/user/cnc_relogin.service.tmpl
[Unit]
Description=Logout then login campus network (05:00 daily via timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

Environment=CNC_USER_ID={{- .env.CNC_USER_ID | quote }}
Environment=CNC_PASSWORD={{- .env.CNC_PASSWORD | quote }}
Environment=CNC_SERVICE={{- .env.CNC_SERVICE | quote }}
ExecStartPre=/home/linuxbrew/.linuxbrew/bin/uvx --from cnc-cli cnc logout
ExecStart=/usr/bin/sleep 2
ExecStartPost=/home/linuxbrew/.linuxbrew/bin/uvx --from cnc-cli cnc login

StandardOutput=journal
StandardError=journal
